# api-monitor.service
[Unit]
Description=API Monitor Service
After=docker.service

[Service]
Type=oneshot
ExecStart=/bin/bash /usr/local/bin/api-monitor.sh

[Install]
WantedBy=multi-user.target
